# taskspine/scripts/systemd/hermes-webhook.service
[Unit]
Description=Hermes Agent (Webhook Gateway)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/hermesagent-taskspine
Environment=HERMES_ENABLE_PROJECT_PLUGINS=1
ExecStart=%h/hermesagent-taskspine/venv/bin/python %h/hermesagent-taskspine/gateway/run.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target

